# nodejd/s6-overlay/etc/services.d/jup/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid root jup -f &>/dev/null
